# systemd/rabbitmq-monitor.service
[Unit]
Description=RabbitMQ Queue Monitor
After=network.target

[Service]
Type=simple
User=rabbitmq-monitor
WorkingDirectory=/opt/rabbitmq-monitor
ExecStart=/opt/rabbitmq-monitor/go-rmq-monitor monitor --config /etc/rabbitmq-monitor/config.yaml
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
